# crontab
* * * * * (date;~/drpeering.com/checkNewCode.bash 2>&1) >> ~/drpeering.com/crontab.log
