# crontab
* * * * * (date;~/drpeering.com/checkNewCode.bash 2>&1) >> /tmp/crontab.log
